# juegos/deploy/juegos.service
[Unit]
Description=Juegos EmulatorJS Go API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/juegos
ExecStart=/opt/juegos/juegos-server -addr 127.0.0.1:8099 -public /opt/juegos/public -emulator /opt/juegos/emulator -roms /opt/juegos/roms -core snes
Restart=on-failure
RestartSec=3
User=www-data
Group=www-data
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
